# distro/overlay/SOPHON/rootfs/etc/systemd/system/sshd-suspend.service
[Unit]
Description=Disconnect SSH before sleep

[Service]
Type=oneshot
ExecStart=/usr/bin/pkill sshd

[Install]
WantedBy=sleep.target

