# distro/overlay/cust01/rootfs/etc/systemd/system/bmrt_setup.service
[Unit]
Description=setup bitmain runtime env.
After=network.target.service

[Service]
User=root
ExecStart=/usr/sbin/bmrt_setup.sh
Type=oneshot

[Install]
WantedBy=multi-user.target
