# ramdisk/target/overlay/bm1684_mix/etc/profile.d/libsophon-bin-path.sh
export PATH=$PATH:/opt/sophon/libsophon-current/bin
export PATH=$PATH:/opt/sophon/libsophon-current/bin/bm1684x
export PATH=$PATH:/opt/sophon/sophon-ffmpeg_0.7.1/bin
export PATH=$PATH:/opt/sophon/sophon-opencv_0.7.1/bin
export PATH=$PATH:/opt/sophon/sophon-sample_0.7.1/bin

export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/opt/sophon/libsophon-current/lib
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/opt/sophon/sophon-ffmpeg_0.7.1/lib
export LD_LIBRARY_PATH=$LD_LIBRARY_PATH:/opt/sophon/sophon-opencv_0.7.1/lib
