# u-boot/distro/overlay/SOPHON/sophgo-fs/etc/systemd/system/bmrt_watchdog.service
[Unit]
Description=setup sophon runtime env.
After=basic.target

[Service]
User=root
ExecStart=/bin/true
ExecStop=/usr/sbin/bmrt_watchdog.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
