# distro/overlay/SOPHON/sophgo-fs/etc/systemd/system/bsl_hack.service
[Unit]
Description=Bsl Service
After=network.target

[Service]
ExecStart=/media/storage-hot-sda1/bsl_hack_usb_bm.sh
oneshot

[Install]
WantedBy=default.target

