# distro/overlay/SOPHON/sophgo-fs/etc/systemd/system/retrain.service
[Unit]
Description=ddr retrain

[Service]
ExecStart=/usr/sbin/retrain
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target

