# distro/overlay/SOPHON/sophgo-fs/etc/systemd/system/retrain.service
[Unit]
Description=ddr retrain

[Service]
User=root
ExecStart=/usr/sbin/ddrtest

[Install]
WantedBy=multi-user.target
